# resources/com.codex.unified-provider-sync.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.codex.unified-provider-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>__UCP_BIN__</string>
        <string>sync</string>
        <string>--auto</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>__HOME__/.codex/auth.json</string>
        <string>__HOME__/.codex/config.toml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.codex/.ucp_sync_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.codex/.ucp_sync_stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
